# 07-metricas/python_flask_app/roles/app/templates/app.conf
[program:app]
command=/usr/bin/python /opt/app/http_server.py
redirect_stderr=true
stdout_logfile=/var/log/app.log
stdout_logfile_maxbytes=100MB
